# udev/fix-greeter-display/90-fix-greeter-display.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/udev-fix-greeter-display.sh"
